# deploy/nhl-game-prediction-ws.service
[Unit]
Description=Webservice backend for prediction mongodb database

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/app/ws
ExecStart=/home/ubuntu/.nvm/versions/node/v13.9.0/bin/node src/server.js
Restart=on-failure
Environment=DATABASE_URL=mongodb://localhost:27017

[Install]
WantedBy=multi-user.target