# deploy/nhl-game-prediction-ws.service
[Unit]
Description=Webservice backend for prediction database

[Service]
Type=simple
User=pi
WorkingDirectory=/app/ws
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
Environment=DATABASE_URL=mongodb://localhost:27017

[Install]
WantedBy=multi-user.target